# misc/upstart.conf
# Put in /etc/upstarat/loopbus.conf
start on x_setup_ready
stop on stopping x

task

script
    /bin/bash /mnt/base-us/scripts/start-sign-from-kindle
end script

